# apps/transcript-server/whisper-server.service
# sudo nvim /etc/systemd/system/whisper-server.service
#
# sudo systemctl daemon-reload
# sudo systemctl enable whisper-server
# sudo systemctl start whisper-server

# View logs
# sudo journalctl -u whisper-server -f
# View status
# sudo systemctl status whisper-server


[Unit]
Description=Whisper.cpp Server
After=network.target

[Service]
Type=simple
User=quolpr
WorkingDirectory=/home/quolpr/whisper.cpp
ExecStart=/home/quolpr/.nvm/versions/node/v22.12.0/bin/node server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=whisper-server

[Install]
WantedBy=multi-user.target
